# templates/node_exporter.service.epp
[Unit]
Description=Node Exporter
After=network.target

[Service]
User=<%= $service_username %>
Group=<%= $service_group %>
Type=simple
ExecStart=<%= $binary_symlink %> --web.config.file=<%= $configuration %><% if $web_listen_address != undef { %> --web.listen-address=<%= $web_listen_address %><% } %>

[Install]
WantedBy=multi-user.target
